# deploy/second-brain.service
[Unit]
Description=second-brain knowledge server
After=network.target

[Service]
Type=simple
# setup.sh が実際のユーザー名に書き換えて配置する
User=%USER%
Group=%USER%
WorkingDirectory=/opt/second-brain
EnvironmentFile=/opt/second-brain/.env
ExecStart=/opt/second-brain/.venv/bin/python -m server.run
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/second-brain

[Install]
WantedBy=multi-user.target
